# deploy/oob-central.service
[Unit]
Description=Out-of-Band Deployment Monitor — Central Service
Documentation=https://github.com/nvlhuongnoi2005/tracking-deployment
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/oob-central --config /etc/oob-central/config.json
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oob-central

[Install]
WantedBy=multi-user.target
